# packages/ufw-geomap/geoip-update.timer
# -*- mode: ini; -*-

[Unit]
Description = Updating GeoIP and GeoLite Databases

[Install]
WantedBy = timers.target

[Timer]
OnCalendar = daily
Persistent = true
