# packages/ufw-geomap/geoip-update.service
# -*- mode: ini; -*-

[Unit]
Description = Updating GeoIP and GeoLite Databases

[Service]
ExecStart = geoipupdate --config-file "%E/GeoIP.conf" --database-directory "%D/geoip/"
