# packages/ufw-geomap/ufw-logs-upload.service
# -*- mode: ini; -*-

[Unit]
Description = Upload UFW Logs to Central Server

[Service]
ExecStart = %h/.local/bin/ufw-logs-upload.sh
