# packages/ufw-geomap/ufw-logs-upload.timer
# -*- mode: ini; -*-

[Unit]
Description = Upload UFW Logs to Central Server

[Install]
WantedBy = timers.target

[Timer]
OnCalendar = daily
Persistent = true
